# docker/cas/supervisord.conf
[supervisord]
nodaemon=true

[program:slapd]
command=/usr/sbin/slapd -u openldap -d 3
redirect_stderr=true
